# syncthing/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/syncthing
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/syncthing
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/stdiscosrv.1.gz
-rw-r--r--	root/root	usr/share/man/man1/strelaysrv.1.gz
-rw-r--r--	root/root	usr/share/man/man1/syncthing.1.gz
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/syncthing-config.5.gz
-rw-r--r--	root/root	usr/share/man/man5/syncthing-stignore.5.gz
drwxr-xr-x	root/root	usr/share/man/man7/
-rw-r--r--	root/root	usr/share/man/man7/syncthing-bep.7.gz
-rw-r--r--	root/root	usr/share/man/man7/syncthing-device-ids.7.gz
-rw-r--r--	root/root	usr/share/man/man7/syncthing-event-api.7.gz
-rw-r--r--	root/root	usr/share/man/man7/syncthing-faq.7.gz
-rw-r--r--	root/root	usr/share/man/man7/syncthing-globaldisco.7.gz
-rw-r--r--	root/root	usr/share/man/man7/syncthing-localdisco.7.gz
-rw-r--r--	root/root	usr/share/man/man7/syncthing-networking.7.gz
-rw-r--r--	root/root	usr/share/man/man7/syncthing-relay.7.gz
-rw-r--r--	root/root	usr/share/man/man7/syncthing-rest-api.7.gz
-rw-r--r--	root/root	usr/share/man/man7/syncthing-security.7.gz
-rw-r--r--	root/root	usr/share/man/man7/syncthing-versioning.7.gz
